# deploy/com.pycronguard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pycronguard</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Library/Frameworks/Python.framework/Versions/3.14/bin/python3</string>
        <string>-m</string>
        <string>pycronguard</string>
        <string>start</string>
        <string>--foreground</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/eden/.pycronguard/logs/launchd_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/eden/.pycronguard/logs/launchd_stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/eden</string>
</dict>
</plist>
